# scripts/linux/systemd/kanban-mcp.service
[Unit]
Description=Kanban MCP Server (user)

[Service]
ExecStart=/usr/bin/env bash -lc 'exec kanban-mcp'
Restart=on-failure

[Install]
WantedBy=default.target
